# systemd/mailqueued.service
[Unit]
Description=Mail queueing daemon

[Service]
ExecStart=/usr/bin/mailqueued

[Install]
WantedBy=default.target
