# systemd/mailqueued.service
[Unit]
After=NetworkManager.service
Requires=NetworkManager.service
Description=Mail queueing daemon

[Service]
ExecStart=/usr/bin/mailqueued

[Install]
WantedBy=multi-user.target
